# packaging/debian/bondster-bco.path
[Unit]
Description=openbank Bondster Marketplace Connection services namespace
PartOf=bondster-bco.service
ConditionPathExists=/etc/bondster-bco/conf.d/init.conf

[Path]
PathChanged=/etc/bondster-bco/conf.d/init.conf
PathModified=/etc/bondster-bco/conf.d/init.conf
Unit=bondster-bco.service

[Install]
WantedBy=bondster-bco.service
