# packaging/debian/bondster-bco-watcher.path
[Unit]
Description=openbank Bondster Marketplace connection configuration watcher
After=bondster-bco.service

[Path]
PathChanged=/etc/bondster-bco/conf.d
Unit=bondster-bco-watcher.service

[Install]
WantedBy=paths.target
